# claude-review-daemon.service
[Unit]
Description=Claude Review Daemon - Polling daemon for Claude Max PR reviews
After=network.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=/opt/claude-review-daemon
ExecStart=/usr/bin/python3 /opt/claude-review-daemon/bridge.py
Restart=always
RestartSec=5

# Graceful shutdown: wait for active reviews to complete
TimeoutStopSec=300
KillSignal=SIGTERM

# Environment
EnvironmentFile=/opt/claude-review-daemon/.env
Environment=CONFIG_FILE=/opt/claude-review-daemon/config.toml

# Ensure Max subscription is used, not API key
Environment=ANTHROPIC_API_KEY=

[Install]
WantedBy=multi-user.target
